# deploy/systemd/tac-master-webhook.service
[Unit]
Description=tac-master webhook server (multi-repo GitHub events)
Documentation=https://github.com/OWNER/tac-master
After=network-online.target tac-master.service
Wants=network-online.target
PartOf=tac-master.service

[Service]
Type=simple
User=krypto
Group=krypto
WorkingDirectory=/srv/tac-master
EnvironmentFile=-/srv/tac-master/config/identity.env
Environment=TAC_MASTER_HOME=/srv/tac-master
Environment=PATH=/usr/local/bin:/usr/bin:/bin:/home/krypto/.local/bin
ExecStart=/home/krypto/.local/bin/uv run orchestrator/webhook_server.py --host 0.0.0.0 --port 8088
Restart=always
RestartSec=10
TimeoutStopSec=30
KillSignal=SIGTERM
KillMode=mixed

NoNewPrivileges=true
ProtectSystem=full
ProtectHome=read-only
ReadWritePaths=/srv/tac-master
PrivateTmp=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictNamespaces=true
LockPersonality=true
RestrictRealtime=true

StandardOutput=append:/srv/tac-master/logs/webhook.stdout.log
StandardError=append:/srv/tac-master/logs/webhook.stderr.log

[Install]
WantedBy=multi-user.target
