# deploy/systemd/tac-master.service
[Unit]
Description=tac-master orchestrator daemon
Documentation=https://github.com/OWNER/tac-master
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=krypto
Group=krypto
WorkingDirectory=/srv/tac-master
EnvironmentFile=-/srv/tac-master/config/identity.env
Environment=TAC_MASTER_HOME=/srv/tac-master
Environment=PATH=/usr/local/bin:/usr/bin:/bin:/home/krypto/.local/bin:/home/krypto/.cargo/bin
ExecStart=/home/krypto/.local/bin/uv run orchestrator/daemon.py
Restart=always
RestartSec=10
# Give the daemon time to finish its current poll cycle on shutdown
TimeoutStopSec=120
KillSignal=SIGTERM
KillMode=mixed

# Security hardening suitable for an LXC container
NoNewPrivileges=true
ProtectSystem=full
ProtectHome=read-only
ReadWritePaths=/srv/tac-master
PrivateTmp=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictNamespaces=true
LockPersonality=true
RestrictRealtime=true

# Logging
StandardOutput=append:/srv/tac-master/logs/daemon.stdout.log
StandardError=append:/srv/tac-master/logs/daemon.stderr.log

[Install]
WantedBy=multi-user.target
