# deploy/systemd/tac-master-dashboard.service
[Unit]
Description=tac-master dashboard (Bun server + WebSocket)
Documentation=https://github.com/kryptobaseddev/tac-master
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=krypto
Group=krypto
WorkingDirectory=/srv/tac-master/dashboard/server
# Load GITHUB_PAT so the /api/github/repo-info endpoint can probe
# private repos the bot has access to
EnvironmentFile=-/srv/tac-master/config/identity.env
Environment=PATH=/usr/local/bin:/usr/bin:/bin:/home/krypto/.bun/bin
Environment=TAC_MASTER_HOME=/srv/tac-master
Environment=TAC_MASTER_DB=/srv/tac-master/state/tac_master.sqlite
Environment=TAC_DASHBOARD_DB=/srv/tac-master/state/dashboard.sqlite
Environment=PORT=4000
ExecStart=/home/krypto/.bun/bin/bun src/index.ts
Restart=always
RestartSec=5
TimeoutStopSec=20

# NoNewPrivileges intentionally OMITTED from the dashboard unit.
# The dashboard exposes POST /api/system/restart-daemon which runs
# `sudo -n systemctl restart tac-master` via a narrow NOPASSWD sudoers
# rule (/etc/sudoers.d/tac-master). NoNewPrivileges=true would block
# that escalation. The sudoers rule is locked to just the three
# tac-master systemctl commands, so the effective blast radius is
# bounded to restarting tac-master's own units — not privilege
# escalation in any meaningful sense.
ProtectSystem=full
ProtectHome=read-only
ReadWritePaths=/srv/tac-master
PrivateTmp=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictNamespaces=true
LockPersonality=true
RestrictRealtime=true

StandardOutput=append:/srv/tac-master/logs/dashboard.stdout.log
StandardError=append:/srv/tac-master/logs/dashboard.stderr.log

[Install]
WantedBy=multi-user.target
